# recipes-support/sysupgrade/files/sysupgrade.service
[Unit]
Description=Data migration

[Service]
ExecStart=/usr/bin/sysupgrade
Type=oneshot
RemainAfterExit=true

[Install]
RequiredBy=basic.target
