# recipes-connectivity/gateway-config-backend/files/gateway-config-backend.socket
[Unit]

[Socket]
ListenStream=80
ListenStream=443

[Install]
WantedBy=sockets.target
Also=gateway-config-backend-sslkey.service
