# recipes-connectivity/firewall/files/firewall.service
[Unit]
Description=install simple iptable rules
Before=network.target
After=network-pre.target

[Service]
Type=oneshot
ExecStart=/usr/lib/seluxit/scripts/firewall

[Install]
WantedBy=multi-user.target
