# recipes-support/led-indicator/files/rf-led.service
[Unit]
Description=RF LED setup
After=ppp.service

[Service]
Type=oneshot
ExecStart=/usr/bin/rf-led-setup

[Install]
WantedBy=multi-user.target

