# recipes-manufacturing/manufacturing-tools/manufacturing-tools/homekit-setup.service
[Unit]
Description=Run HomeKit setup during manufacturing
Requires=ipr.service manufacturing-statusfiles.service
After=ipr.service manufacturing-statusfiles.service
Before=serial-getty@ttyS0.service
ConditionPathExists=!/etc/hk_setup_done

[Service]
ExecStart=/usr/bin/homekit-setup
Type=oneshot
ExecStopPost=/bin/sh -c 'if [ "$$EXIT_STATUS" != 0 ]; then fct-tool --set-leds red; fi'

[Install]
WantedBy=serial-getty@ttyS0.service
