# recipes-manufacturing/manufacturing-tools/manufacturing-tools/ipr.service
[Unit]
Description=Run IPR tool
Before=serial-getty@ttyS0.service

[Service]
ExecStart=/usr/bin/ipr-setup
Type=oneshot
