# recipes-manufacturing/manufacturing-tools/manufacturing-tools/ipr.service
[Unit]
Description=Run IPR tool
Requires=manufacturing-statusfiles.service
After=manufacturing-statusfiles.service
Before=serial-getty@ttyS0.service
ConditionPathExists=!/etc/ipr_setup_done

[Service]
ExecStart=/usr/bin/ipr-setup
Type=oneshot
RemainAfterExit=yes
ExecStopPost=/bin/sh -c 'if [ "$$EXIT_STATUS" != 0 ]; then cpms-errorhandler; fi'

[Install]
WantedBy=serial-getty@ttyS0.service
